# conf/macos/extra-files/Library/LaunchDaemons/org.openobservatory.ooniprobe-agent.plist
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
 <key>Label</key><string>org.openobservatory.ooniprobe-agent.startup</string>
 <key>EnvironmentVariables</key>
 <dict>
   <key>PYTHONPATH</key>
   <string>/System/Library/Frameworks/Python.framework/Versions/2.7/lib/python2.7</string>
 </dict>
 <key>ProgramArguments</key>
 <array>
    <string>/Library/Application Support/ooniprobe/venv/ooniprobe/bin/python</string>
    <string>/Library/Application Support/ooniprobe/venv/ooniprobe/bin/ooniprobe-agent</string>
    <string>run</string>
 </array>
 <key>RunAtLoad</key>
   <true/>
 <key>KeepAlive</key>
   <true/>
 <key>StandardErrorPath</key>
   <string>/dev/null</string>
 <key>StandardOutPath</key>
   <string>/dev/null</string>
</dict>
</plist>
