# conf/pacman/extra-files/etc/init/ooniprobe.conf
description "ooniprobe-agent system service"

start on runlevel [2345]
stop on runlevel [016]

env HOME="/var/lib/ooni"
exec sudo -u arch-ooni /usr/bin/ooniprobe-agent start
